# gunicorn.service
[Unit]
Description=demo fastapi application
After=network.target

[Service]
User=bhdima
Group=bhdima
WorkingDirectory=/home/bhdima/app/src/
Environment="PATH=/home/bhdima/app/venv/bin"
EnvironmentFile=/home/bhdima/.env
ExecStart=/home/bhdima/app/venv/bin/gunicorn -w 4 -k uvicorn.workers.UvicornWorker app.main:app --bind 0.0.0.0:8000

[Install]
WantedBy=multi-user.target